# systemfiles/rc.local
# This file initialises everything
# It will add the cronjob that begins
# to handle everything.

# before the `exit 0` line just inject this:


if ! crontab -upi -l | grep "run.sh" >/dev/null; then
	crontab -upi /boot/picluster/picron
fi